# systemd/wspr-uploader.service
[Unit]
Description=wspr-recorder standalone uploader (drains sink.db → wsprnet/wsprdaemon/pskreporter)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=notify
User=wsprrec
Group=wsprrec

WatchdogSec=180
NotifyAccess=main

EnvironmentFile=-/etc/sigmond/coordination.env
# NOT under env/ — that dir is scanned for wspr-recorder@<instance> decoders,
# and an env/uploader.env there spawns a phantom wspr-recorder@uploader unit.
EnvironmentFile=-/etc/wspr-recorder/uploader.env

WorkingDirectory=/opt/git/sigmond/wspr-recorder

# Runs ONLY the hs-uploader pump — no radiod, no decode, no config.toml.
# Identity / SFTP / merge / verify / audit all come from the env file above
# via WsprUploaderHs.from_env().
ExecStart=/opt/git/sigmond/wspr-recorder/venv/bin/python3 -m wspr_recorder.cli uploader

# No radiod to provision — the pump is up in well under a second.
TimeoutStartSec=30

Restart=always
RestartSec=5

StandardOutput=journal
StandardError=inherit
SyslogIdentifier=wspr-uploader

Nice=5
MemoryMax=512M
MemorySwapMax=0

NoNewPrivileges=true
# Private writable /tmp: hs-uploader's _upload_tar stages the .tbz via
# tempfile.NamedTemporaryFile (default tmpdir), and ProtectSystem=strict
# makes the host /tmp read-only.  The uploader has no host-/tmp dependency
# (it talks to decoders via sink.db + upload-wake.sock under /var/lib/sigmond),
# so an isolated private /tmp is the clean fix.
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
# /var/lib/sigmond: the shared sink.db (pending_uploads) + upload-wake.sock.
# /var/lib/hs-uploader: the uploader's watermark store.
# No /dev/shm/wspr-recorder, no StateDirectory (callhash) — those are
# decode-side only.
ReadWritePaths=/var/lib/sigmond /var/lib/hs-uploader
ReadOnlyPaths=/etc/wspr-recorder
PrivateNetwork=false
# Outbound SFTP/HTTP only — no raw-socket / multicast RTP capture here, so
# none of the recorder's CAP_NET_RAW / CAP_NET_BIND_SERVICE are needed.

[Install]
WantedBy=multi-user.target
